# blobe-optimizer.service
[Unit]
Description=Blobe VM Optimizer Service
After=network.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/blobe-vm
# Call the idempotent ensure script which will install prerequisites
# and then exec the node process so systemd tracks the optimizer PID.
ExecStart=/usr/bin/env bash /opt/blobe-vm/optimizer/optimizer-ensure.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
